# homedir/.crontab
# ruby gem nightly updates
0 4 * * *  gem update --verbose
# brew nightly updates
10 4 * * *  brew update
20 4 * * *  brew upgrade
# brew cask
30 4 * * *  brew cask upgrade
# mas updates
40 4 * * *  mas upgrade
# Others
50 4 * * *  tldr --update
